# root/etc/services.d/cardigann/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc "${RUN_OPTS}" cardigann server
